# deploy/agent/deuswatch-agent.service
[Unit]
Description=DeusWatch Agent
Documentation=https://github.com/FirdausYudha/DeusWatch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/deuswatch/agent.env
ExecStart=/usr/local/bin/deuswatch-agent
Restart=always
RestartSec=5
# Read /var/log read-only without full root (hardening).
DynamicUser=yes
SupplementaryGroups=adm systemd-journal
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/deuswatch

[Install]
WantedBy=multi-user.target
